# Git-My-Way-MCP/deploy/git-my-way-mcp.service
[Unit]
Description=Git-My-Way-MCP (Scaffold)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 %h/Projects/Reusable-MCP/Git-My-Way-MCP/server/app.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target

