# SysDiag-MCP/deploy/sysdiag-mcp.service
[Unit]
Description=SysDiag-MCP (Scaffold)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 %h/Projects/Reusable-MCP/SysDiag-MCP/server/app.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target

